# Link-Up-for-impact-main/uploaded/LinkUp-main/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# ============================================
# Nginx - Servidor Web para Frontend
# ============================================
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=10

# ============================================
# Uvicorn - Servidor FastAPI Backend
# ============================================
[program:uvicorn]
command=uvicorn server:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app/backend
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/uvicorn.err.log
stdout_logfile=/var/log/supervisor/uvicorn.out.log
priority=20
environment=PYTHONUNBUFFERED="1"
